# crontab
CRON_TZ=America/Phoenix

30 21 * * * bun run etl