# crontab
CRON_TZ=America/Phoenix

30 21 * * * yarn run etl